# desktop-portal/xdg-desktop-portal-thedesk.service.in
[Unit]
Description=XDG Portal (theDesk)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.thedesk
ExecStart=${CMAKE_INSTALL_FULL_LIBDIR}/thedesk-desktop-portal
Restart=no
